# ssh.conf
[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/data/logs/ssh.err.log
stdout_logfile=/data/logs/ssh.out.log